# base/conf/shell/etc/profile.d/00-reset-path.sh
# Reset default PATH settings
PATH="/opt/code-server/bin:$JULIA_PATH/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
export PATH
